# netifd-morse/etc/hotplug.d/ieee80211/01-remove-auto-vif
#!/bin/sh

# In general, drivers automatically bring up a single managed vif at wlan?
# (see AUTO_VIF).
#
# In our case, since we want netifd to manage everything, we don't want this. In
# particular, it's a problem if wlan0 comes up when the netifd wifi is completely
# disabled. This happens on first boot with our EKH devices, as they're
# disabled by default since we don't know the user's region.
#
# But the real reason we need this is to make sure we can use the wlan?
# numbering for the _morse_ devices, rather than have some other
# mac80211 device steal it (once mac80211.sh is loaded by netifd,
# they'll use phy?-role? naming). So this 'belongs' in netifd-morse
# because it's effectively part of our morse.sh setup, even though
# it applies to all wifi devices.
for device in /sys$DEVPATH/device/net/wlan*; do
	iw dev "$(basename "$device")" del
done
