# utils/antenna-select/files/etc/init.d/antenna-select
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2021-2022 Morse Micro Pty Ltd. All rights reserved.
#

# After boot for overlay mount, (and log).
# Before wpad and network.
START=15

boot() {
    /sbin/antenna-select 2
}
